# conf/cf/etc/cron.daily/hostcheck
#! /bin/sh
source /usr/local/sbin/syscheck
